# z_docker/etc/s6-overlay/s6-rc.d/20-web/run
#!/command/with-contenv bash

cd /app
exec npm run web